# launchd/com.tanmayrsm.dlyprb.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tanmayrsm.dlyprb.autopush</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-lc</string>
    <string>/usr/bin/osascript -e 'tell application "Terminal" to do script "REPO_DIR=/Users/tamishra/Documents/dlyPrb BRANCH=main GITHUB_USERNAME=tanmayrsm GITHUB_TOKEN_FILE=/Users/tamishra/.config/dlyPrb/github_token IST_TZ=Asia/Kolkata RUN_WINDOW_HOURS=10,12 STATE_DIR=/Users/tamishra/.local/state/dlyPrb-autopush /Users/tamishra/Documents/dlyPrb/scripts/auto-push-dlyprb.sh &gt;&gt; /tmp/dlyprb-autopush.out.log 2&gt;&gt; /tmp/dlyprb-autopush.err.log; exit"'</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>REPO_DIR</key>
    <string>/Users/tamishra/Documents/dlyPrb</string>
    <key>BRANCH</key>
    <string>main</string>
    <key>GITHUB_USERNAME</key>
    <string>tanmayrsm</string>
    <key>GITHUB_TOKEN_FILE</key>
    <string>/Users/tamishra/.config/dlyPrb/github_token</string>
    <key>IST_TZ</key>
    <string>Asia/Kolkata</string>
    <key>RUN_WINDOW_HOURS</key>
    <string>10,12</string>
    <key>STATE_DIR</key>
    <string>/Users/tamishra/.local/state/dlyPrb-autopush</string>
  </dict>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/tmp/dlyprb-autopush.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/dlyprb-autopush.err.log</string>
</dict>
</plist>
